# docs/config_file_sample/domain_supervisor.conf
# -*- conf -*-
[include]
files = *.supervisor

[supervisord]

[supervisorctl]
serverurl = unix://supervisord.sock

[unix_http_server]
file = supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:main]
process_name = main-%(process_num)s
command = python /home/tuziidev/dev_blog2/blog/runserver.py --port=%(process_num)s
numprocs = 4
numprocs_start = 8000
directory=/home/tuziidev/
redirect_stderr=true
stdout_logfile=/home/tuziidev/log/tuzii_main%(process_num)s
